# apkovl/etc/init.d/megabit
#!/sbin/openrc-run
name=megabit
command="/usr/local/bin/megabit"
pidfile="/var/run/megabit.pid"
