# apkovl/etc/init.d/megabit-sync
#!/sbin/openrc-run
name=megabit-sync
command="/usr/local/bin/megabit-sync"
pidfile="/var/run/megabit-sync.pid"

start_pre() {
  mkdir -p /mnt/tmp
  mount /usr/share/megabit/tmp.img /mnt/tmp
  modprobe -r g_acm_ms
  modprobe g_acm_ms file=/usr/share/megabit/tmp.img
}
