# apkovl/etc/inittab
::sysinit:/sbin/openrc sysinit
::sysinit:/sbin/openrc boot
::wait:/sbin/openrc default
ttyGS0::respawn:/sbin/getty -L ttyGS0 115200 vt100
